# src/main/assembly/dist/bin/easy-sword2.service
[Unit]
Description=EASY Sword2 Service

[Service]
ExecStart=/opt/dans.knaw.nl/easy-sword2/bin/easy-sword2-service-start

User=easy-sword2
Group=easy-sword2

[Install]
WantedBy=multi-user.target
